# supervisor-configs/scheduler.conf
[program:scheduler]
command=sh -c '/exe/kube-scheduler --master=$KUBE_API_ADDRESS:$KUBE_API_PORT $SCHEDULER_ARGS'
stdout_events_enabled = true
stderr_events_enabled = true