# supervisor-configs/scheduler.conf
[program:scheduler]
command=sh -c '/exe/kube-scheduler --master=127.0.0.1:$KUBE_API_PORT $SCHEDULER_ARGS'
stdout_events_enabled = true
stderr_events_enabled = true
